# udp_connection/udp_to_drive/udp_to_drive.service
[Unit]
Description=SAFRS UDP → ROS2 Drive Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ros2_ws
Environment="PYTHONUNBUFFERED=1"

# === ROS2 Environment ===
ExecStart=/bin/bash -c "\
    source /opt/ros/humble/setup.bash && \
    source /home/ubuntu/ros2_ws/install/setup.bash && \
    ros2 launch udp_to_drive udp_bridge.launch.py \
"

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
